# docker/admin_server/supervisor_scrapyd.ini
[program:scrapyd]
command=scrapyd
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/scrapyd/scrapyd.out.log
stderr_logfile=/var/log/scrapyd/scrapyd.err.log

[program:scraper_scheduler]
command=python3 /scraper/scraper_scheduler.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/scraper_scheduler/scraper_scheduler.out.log
stderr_logfile=/var/log/scraper_scheduler/scraper_scheduler.err.log
