# docker/admin_server/supervisor_scrapyd.ini
[program:scrapyd]
command=scrapyd --pidfile ./scrapy.pid
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/scrapyd/stdout.log
stderr_logfile=/var/log/scrapyd/stderr.log

[program:scraper_scheduler]
command=python3 /scraper/scraper_scheduler.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/scraper_scheduler/stdout.log
stderr_logfile=/var/log/scraper_scheduler/stderr.log
