# docker/admin_server/supervisor_scrapyd.ini
[program:scrapyd]
command=scrapyd
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/scrapyd/scrapyd.err.log
stdout_logfile=/var/log/scrapyd/scrapyd.out.log
